# munus.service
[Unit]
Description=Munus — volunteer hour tracker
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/munus
EnvironmentFile=/home/pi/munus/.env
ExecStart=/home/pi/munus/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8001
Restart=on-failure

[Install]
WantedBy=multi-user.target
